# roles/vllm/templates/vllm-sanitizer.service.j2
{{ ansible_managed | comment }}
# Image-sanitizing reverse proxy that sits in front of vLLM. Catches
# the FB CDN URLs whose bytes Pillow can't decode (animated WebP, signed
# 403 HTML, SVG data URIs, etc.) before they reach vLLM, so the model
# never sees the crash-inducing payload.

[Unit]
Description=vLLM image-sanitizing proxy (port {{ vllm_public_port }} → 127.0.0.1:{{ vllm_internal_port }})
After=vllm.service
Requires=vllm.service

[Service]
Type=simple
User={{ vllm_user }}
Group={{ vllm_user }}
WorkingDirectory={{ vllm_workdir }}
Environment="VLLM_UPSTREAM=http://127.0.0.1:{{ vllm_internal_port }}"
Environment="SANITIZER_HOST={{ vllm_public_host }}"
Environment="SANITIZER_PORT={{ vllm_public_port }}"
Environment="SANITIZER_FETCH_TIMEOUT={{ vllm_sanitizer_fetch_timeout }}"
Environment="SANITIZER_JPEG_QUALITY={{ vllm_sanitizer_jpeg_quality }}"
Environment="SANITIZER_MAX_EDGE={{ vllm_sanitizer_max_edge }}"
ExecStart={{ vllm_workdir }}/.venv/bin/python {{ vllm_workdir }}/image_sanitizer.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
