# roles/observability/templates/ollama-exporter.service.j2
{{ ansible_managed | comment }}
# Rendered from roles/observability/templates/ollama-exporter.service.j2.
# frcooper/ollama-exporter polls /api/ps and /api/tags. Live-state
# signal (warm models, RAM per model). Canary handles regression.

[Unit]
Description=frcooper/ollama-exporter (Ollama live-state metrics)
After=network-online.target ollama.service
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/ollama-exporter \
    --ollama-url=http://127.0.0.1:11434 \
    --listen-address=:{{ obs_ollama_exporter_port }}
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
